# scripts/deploq.service
[Unit]
Description=Deploq — Lightweight Git Push Deploy
After=network.target docker.service
Requires=docker.service

[Service]
Type=simple
User=deploy
Group=deploy
WorkingDirectory=/opt/deploq
ExecStart=/usr/local/bin/deploq serve -config /etc/deploq/deploq.yaml -log-format json
Restart=on-failure
RestartSec=5s
EnvironmentFile=/etc/deploq/env

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/home /opt /var/log

[Install]
WantedBy=multi-user.target
