# scripts/pushup.service
[Unit]
Description=Pushup — Lightweight Git Push Deploy
After=network.target docker.service
Requires=docker.service

[Service]
Type=simple
User=deploy
Group=deploy
WorkingDirectory=/opt/pushup
ExecStart=/usr/local/bin/pushup serve -config /etc/pushup/pushup.yaml -log-format json
Restart=on-failure
RestartSec=5s
EnvironmentFile=/etc/pushup/env

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/home /opt /var/log

[Install]
WantedBy=multi-user.target
